# .config/systemd/user/mbsync.service
[Unit]
Description=Maildir syncing

[Service]
Type=oneshot
ExecStart=mbsync -Va
