# .config/systemd/user/notmuch.service
[Unit]
Description=Notmuch indexing service

[Service]
Type=oneshot
ExecStart=/usr/bin/notmuch new
